# etc/cron.d/resolve_infra
* * * * * /usr/local/bin/reconfig-dns.sh /etc/dnsmasq.conf.template > /etc/dnsmasq.conf && systemctl restart dnsmasq
